# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:hamper]
command=/hamper/scripts/hamper
environment=PYTHONPATH="."
directory=/hamper
autostart=true
autorestart=true

